# sys/apps/openrc/extras/disable-secondary-gpu.initd
#!/sbin/openrc-run

name=disable-secondary-gpu

depend() {
    before dbus
}

start() {
    /usr/libexec/disable-secondary-gpu
}
