# sys/power/upower/upowerd.initd
#!/sbin/openrc-run
# Copyright (c) 2008 Saleem Abdulrasool <compnerd@compnerd.org>
# Distributed under the terms of the GNU General Purpose License v2
# $Header: $

description="Power Manager Daemon"

start(){
    /usr/libexec/upowerd &
    return 0
}
stop(){
    return 0
}

depend() {
   need dbus
}

# vim: set ft=gentoo-init-d ts=3 sw=3 et:
